# vm/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 %(ENV_RESOLUTION)sx24
autorestart=true
priority=100

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -xkb -forever -shared -repeat -capslock -passwd cyberlab -rfbport 5901
autorestart=true
priority=200

[program:novnc]
command=/usr/bin/websockify --web /usr/share/novnc 6080 localhost:5901
autorestart=true
priority=300

[program:xfce4]
command=/usr/bin/startxfce4
environment=DISPLAY=":1",HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s"
user=%(ENV_USER)s
autorestart=true
priority=400
